# debian/kurz.postinst
#!/bin/sh

systemctl enable /lib/systemd/system/kurz.target
systemctl start kurz.target
