# debian/kurz.prerm
#!/bin/sh

systemctl stop kurz.target
systemctl disable /lib/systemd/system/kurz.target
