# src/init/kurz.target
[Unit]
Description=Kurz jabber-bot
After=multi-user.target
